# clab-frr/frr-start.service
[Unit]
Start FRR service
After=network-online.target

[Service]
ExecStart=/etc/systemd/system/frr-start-svc.sh

[Install]
WantedBy=multi-user.target
